# systemd/user/mouse-gesture.service
[Unit]
Description=Magic Mouse Gesture Recognition
Documentation=https://github.com/jiahaoxiang2000/mouse-gesture
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mouse-gesture-recognition --config /home/isomo/.config/mouse-gesture/config.json
Restart=always
RestartSec=5
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
